# deploy/systemd/qikuku-billing-monthly.timer
[Unit]
Description=Run Qikuku monthly billing grant daily at 02:00 UTC

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true
Unit=qikuku-billing-monthly.service

[Install]
WantedBy=timers.target
